# tooling/docker/guest-runner-init-script
#!/sbin/openrc-run
command="/usr/bin/guest-runner"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"